# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
loglevel=info
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Asterisk SIP bridge (registers ext 101 with Alpha PBX) ───────────────────
[program:asterisk]
command=asterisk -f
autostart=true
autorestart=true
stderr_logfile=/tmp/asterisk.err.log
stdout_logfile=/tmp/asterisk.out.log

; ── AudioSocket bridge (Asterisk ↔ Pipecat WebSocket) ─────────────────────────
[program:bridge]
command=python /app/bridge_to_ws.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/bridge.err.log
stdout_logfile=/tmp/bridge.out.log

; ── Pipecat FastAPI server ─────────────────────────────────────────────────────
[program:pipecat]
command=python -m uvicorn server:app --host 0.0.0.0 --port 7860
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/pipecat.err.log
stdout_logfile=/tmp/pipecat.out.log
